# aiproxy/src/aiproxy/share/com.github.olafgeibig.tools.aiproxy.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>com.github.olafgeibig.tools.aiproxy</string>
	<key>ProgramArguments</key>
	<array>
		<string>{{AIPROXY_PATH}}</string>
	</array>
	<key>RunAtLoad</key>
	<true/>
	<key>KeepAlive</key>
	<true/>
	<key>StandardOutPath</key>
	<string>{{STDOUT_PATH}}</string>
	<key>StandardErrorPath</key>
	<string>{{STDERR_PATH}}</string>
</dict>
</plist>
                                                                                  